# ft8900r-watchdog.service
[Service]
ExecStart=/usr/local/bin/ft8900r-watchdog
Restart=always

[Install]
WantedBy=multi-user.target
